# templates/systemd/app-watcher.path
[Path]
PathChanged=/var/app

[Install]
WantedBy=multi-user.target
