# templates/systemd/config-monitor.path
[Unit]
Description=Restart app on config changes
After=app.service

[Path]
PathChanged=/var/app/.env
PathChanged=/var/app/docker-compose.yaml
Unit=config-monitor.service

[Install]
WantedBy=multi-user.target
